# cribsense.service.in
[Unit]
Description=Cribsense Baby Monitoring Service

[Service]
Type=simple
ExecStart=@bindir@/cribsense --config @sysconfdir@/cribsense/config.ini
CapabilityBoundingSet=
ProtectSystem=true
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
